# transmission-idle-shutdown.service
[Unit]
Description=Clean up finished torrents and check for idle shutdown
After=network.target

[Service]
Type=oneshot
User=ervin
ExecStart=/home/ervin/.local/bin/transmission-watch.sh
