# transmission-idle-shutdown.service
[Unit]
Description=Clean up finished torrents and check for idle shutdown
After=network.target transmission-daemon.service

[Service]
Type=oneshot
User=ervin
ExecStart=/home/ervin/.local/bin/rsync-torrents watch
Environment=PATH=/home/ervin/.local/bin:/usr/bin:/bin
ProtectSystem=full
PrivateTmp=true
NoNewPrivileges=true
ProtectControlGroups=true
ProtectKernelModules=true
ProtectKernelTunables=true
RestrictRealtime=true
StandardOutput=journal
StandardError=journal
