# transmission-proxy.socket
[Unit]
Description=Transmission RPC activation socket

[Socket]
ListenStream=9091
Accept=no

[Install]
WantedBy=sockets.target
