# transmission-idle-shutdown.service
[Unit]
Description=Clean up finished torrents and check for idle shutdown
After=network.target transmission-daemon.service
Wants=transmission-daemon.service

[Service]
Type=oneshot
User=@USER@
ExecStart=@BINDIR@/rsync-torrents watch
Environment=PATH=@BINDIR@:/usr/bin:/bin
TimeoutStartSec=120
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
NoNewPrivileges=true
ProtectControlGroups=true
ProtectKernelModules=true
ProtectKernelTunables=true
RestrictRealtime=true
StandardOutput=journal
StandardError=journal
