# transmission-idle-shutdown.service
[Unit]
Description=Clean up finished torrents and check for idle shutdown
After=network.target transmission-daemon.service

[Service]
Type=oneshot
ExecStart=/home/ervin/.local/bin/rsync-torrents -c /home/ervin/.config/rsync-torrents/config.toml watch
Environment=PATH=/home/ervin/.local/bin:/usr/bin:/bin
Environment=HOME=/home/ervin
Environment=XDG_RUNTIME_DIR=/run
StandardOutput=journal
StandardError=journal
